# package/factory-defaults/files/etc/hotplug.d/button/99-reset-to-factory-defaults
#!/bin/sh

if [ "$BUTTON" = reset -a "$ACTION" = pressed ]; then
	firstboot                      
	reboot                      
fi

exit 0
